# deploy/pocketbase.service
[Unit]
Description=PocketBase - TCN Coffee Management Backend
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/pocketbase
ExecStart=/opt/pocketbase/pocketbase serve --http=127.0.0.1:8090
Restart=always
RestartSec=5
StandardOutput=append:/var/log/pocketbase.log
StandardError=append:/var/log/pocketbase.log

[Install]
WantedBy=multi-user.target
